# archiso/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
# BlazzCore — Auto-login on TTY1
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin blazzcore --noclear %I $TERM
